# root/etc/systemd/system/snapper-cleanup.service.d/mutex.conf
[Service]
Type=oneshot
ExecStart=
ExecStart=/usr/bin/flock -x -F /media/auxiliary/.snapshot_hist.lock     \
    /usr/lib/snapper/systemd-helper --cleanup
