# root/etc/systemd/system/snapper-cleanup.service.d/mutex.conf
[Unit]
Before=upload-snapshots.service archive-snapshots.service